# overlay/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin flavos --noclear %I $TERM
Type=idle
