# overlay/etc/profile.d/flavos.sh
# Flavos OS Global Environment Variables
export GTK_CSD=0
export NO_AT_BRIDGE=1
